# adguard-home/adguard-ntfy-report.timer
[Unit]
Description=AdGuard Home — rapport quotidien ntfy à 8h

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
